# images/os/usr/lib/systemd/system/k3s.service.d/override.conf
[Service]
ConfigurationDirectory=rancher/k3s rancher/node
ConfigurationDirectoryMode=0700
LogsDirectory=rancher/k3s
LogsDirectoryMode=0700
StateDirectory=rancher/k3s
StateDirectoryMode=0700

ExecStart=
ExecStart=/usr/local/bin/k3s server --config=/usr/lib/rancher/k3s/config.yaml
